# .config/systemd/user/power-profile-monitor.timer
[Unit]
Description=Run power profile check every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
