# .config/systemd/user/power-profile-monitor.service
[Unit]
Description=Check power profile thresholds

[Service]
Type=oneshot
ExecStart=/usr/local/bin/power-profile-auto.sh
